# debian/ryu-bin.postrm
update-rc.d -f ryu remove >/dev/null || exit $?
#DEBHELPER#
